# docs/setup/saas_celery_supervisor.conf
# NexusCore SaaS - Celery Worker supervisor 設定例
#
# 配置場所: /etc/supervisor/conf.d/nexuscore-celery.conf
#
# 使用方法:
#   sudo cp docs/saas_celery_supervisor.conf /etc/supervisor/conf.d/nexuscore-celery.conf
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start nexuscore-celery:*
#
# ログ確認:
#   sudo tail -f /var/log/nexuscore/celery-worker.log

[program:nexuscore-celery]
command=/home/nexuscore/NexusCore/myenv_linux/bin/celery -A nexuscore.webapp.celery_app.celery worker --loglevel=INFO --concurrency=4
directory=/home/nexuscore/NexusCore
user=nexuscore
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nexuscore/celery-worker.log
environment=FLASK_APP="nexuscore.webapp:create_app",CELERY_BROKER_URL="redis://localhost:6379/0",CELERY_RESULT_BACKEND="redis://localhost:6379/1",NEXUS_USE_CELERY="1"

