# docs/setup/saas_celery_systemd.service
# NexusCore SaaS - Celery Worker systemd サービス設定例
#
# 配置場所: /etc/systemd/system/nexuscore-celery.service
#
# 使用方法:
#   sudo cp docs/saas_celery_systemd.service /etc/systemd/system/nexuscore-celery.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable nexuscore-celery
#   sudo systemctl start nexuscore-celery
#
# ログ確認:
#   sudo journalctl -u nexuscore-celery -f

[Unit]
Description=NexusCore Celery Worker
After=network.target redis.service

[Service]
Type=simple
User=nexuscore
Group=nexuscore
WorkingDirectory=/home/nexuscore/NexusCore
Environment="PATH=/home/nexuscore/NexusCore/myenv_linux/bin"
Environment="FLASK_APP=nexuscore.webapp:create_app"
Environment="CELERY_BROKER_URL=redis://localhost:6379/0"
Environment="CELERY_RESULT_BACKEND=redis://localhost:6379/1"
Environment="NEXUS_USE_CELERY=1"
ExecStart=/home/nexuscore/NexusCore/myenv_linux/bin/celery -A nexuscore.webapp.celery_app.celery worker --loglevel=INFO --concurrency=4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

